# supervisor.conf
[supervisord]
nodaemon=true

[program:hbase_master]
command=/etc/init.d/hbase-master start
priority=10

[program:hbase_thrift]
command=/etc/init.d/hbase-thrift start
priority=100
